# docker/images/php/container-files/etc/supervisor.d/install_wordpress.conf
[program:install-wordpress]
directory = /
command = /bin/bash install-wp.sh
autostart=true
autorestart=false
startsecs=10
startretries=0
stopsignal=TERM
stopwaitsecs=5